# palfynux-pkg/palfynux.install
post_install() {
  systemctl enable lightdm;
  systemctl start;
  echo "Palfynux telepítve!"
}
