# palfynux-pkg/palfynux.install
post_install() {
  update-mime-database /usr/share/mime;
  update-desktop-database;
  systemctl enable sshd;
  systemctl start sshd;
  systemctl enable lightdm;
  systemctl start lightdm;
  echo "Palfynux telepítve!"
}
